# templates/default/openshift.service.erb
[Unit]
Description=Openshift

[Service]
ExecStart=/usr/local/bin/openshift start

[Install]
WantedBy=multi-user.target
